# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
pidfile=/var/run/supervisord.pid
loglevel=info

[program:api]
command=uvicorn src.api.main:app --host %(ENV_API_HOST)s --port %(ENV_API_PORT)s
directory=/app
autostart=true
autorestart=unexpected
startretries=3
startsecs=5
stopwaitsecs=10
stopasgroup=true
killasgroup=true
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:web]
command=node /app/frontend/server.js
directory=/app
autostart=true
autorestart=unexpected
startretries=3
startsecs=3
stopwaitsecs=5
stopasgroup=true
killasgroup=true
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
